# dotfiles/systemd/low-bat.service
[Unit]
Description=Suspend laptop at low battery.

[Service]
ExecStart=/home/yobibyte/scripts/bat.sh

