# dotfiles/systemd/low-bat.timer
[Unit]
Description=Suspend laptop at low battery.

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target

