# dotfiles/systemd/newsboat.timer
[Unit]
Description=Load new newsboat entries.

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target

