# dotfiles/systemd/newsboat.service
[Unit]
Description=Load new newsboat entries.

[Service]
ExecStart=/home/yobibyte/scripts/reload_newsboat.sh

